# deploy/aarsp-sidecar.service
# Copyright 2025 KalGuard Contributors. Licensed under the Apache License, Version 2.0.
# Systemd unit for KalGuard sidecar. Install to /etc/systemd/system/
# Then: systemctl daemon-reload && systemctl enable --now kalguard-sidecar
[Unit]
Description=KalGuard AI Agent Runtime Security Platform Sidecar
Documentation=https://github.com/kalguard/kalguard
After=network.target

[Service]
Type=simple
User=kalguard
Group=kalguard
WorkingDirectory=/opt/kalguard
ExecStart=/usr/bin/node dist/cmd/sidecar.js
Restart=on-failure
RestartSec=5
EnvironmentFile=-/etc/kalguard/env
Environment=KALGUARD_HOST=0.0.0.0
Environment=KALGUARD_PORT=9292
Environment=KALGUARD_POLICY_PATH=/etc/kalguard/policy.json
Environment=KALGUARD_AUDIT_LOG_PATH=/var/log/kalguard/audit.log
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=read-only
ReadWritePaths=/var/log/kalguard
ReadOnlyPaths=/opt/kalguard /etc/kalguard

[Install]
WantedBy=multi-user.target
